# multistagebuild/etc/systemd/system/snap-snapd-17883.mount
[Unit]
Description=Mount unit for snapd, revision 17883
Before=snapd.service
After=zfs-mount.service

[Mount]
What=/var/lib/snapd/snaps/snapd_17883.snap
Where=/snap/snapd/17883
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=multi-user.target
